# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
autostart=true
autorestart=true

[program:php]
command=/usr/sbin/php -D FOREGROUND
autostart=true
autorestart=true